# templates/default/gmetric/mongodb.cron.erb
* * * * * root MONGODB_PORT=27018 /usr/bin/python /usr/local/bin/mongodb-ganglia
