# templates/default/gmetric/tcp-connections.cron.erb
* * * * * root /usr/bin/perl /usr/local/bin/tcp-connections-ganglia