# templates/default/gmetric/tcp-connections.cron.erb
* * * * * root /usr/local/bin/tcp-connections-ganglia